# rfid_cam/systemd/rfid_cam.service
[Unit]
Description=RFID Pet Chip Reader with Immediate Notifications
After=network-online.target dev-ttyUSB1.device
Wants=network-online.target

[Service]
Type=simple
User=collins
Group=collins
WorkingDirectory=/home/collins/repos/pet-chip-reader/rfid_cam
ExecStart=/usr/bin/python3 /home/collins/repos/pet-chip-reader/rfid_cam/src/single_camera_test.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

# Groups for hardware access
SupplementaryGroups=video dialout gpio

# Security settings (relaxed for camera/serial access)
NoNewPrivileges=true
ProtectSystem=false
ProtectHome=false
PrivateTmp=false

# Resource limits
MemoryMax=512M
CPUQuota=80%

[Install]
WantedBy=multi-user.target